# services/rgpio_monitor/run
#!/bin/sh
echo "*** starting WSRemoteGPIO ***"
exec 2>&1
exec /data/WSRemoteGPIO/rgpio_monitor
